# scan.supervisor.conf
[program:scan]
command=/usr/bin/python /home/pi/scan/scan.py
directory=/home/pi/scan
redirect_stderr=true
startsecs=5
autorestart=true
stdout_logfile=/home/pi/scan/main.stdout.log
